# 01-soc/rules/wazuh_ad_attacks.xml
<!-- 
  Wazuh Custom Rules — Active Directory Attack Detection
  Author: Joe Bichall (@joebat10)
  Référence: MITRE ATT&CK Enterprise Matrix
  Source: Adaptées du lab "Attacking Active Directory with Linux"
  
  Basé sur les Event IDs Sysmon + Windows Security
  Deployment: /var/ossec/etc/rules/wazuh_ad_attacks.xml
-->
<group name="ad_attacks,windows,">

  <!-- ============================================================ -->
  <!-- T1003.001 — LSASS Memory Access (Credential Dumping)        -->
  <!-- Sysmon Event ID 10 : ProcessAccess                          -->
  <!-- ============================================================ -->
  <rule id="100001" level="15">
    <if_group>sysmon_event10</if_group>
    <field name="win.eventdata.targetImage" type="pcre2">(?i)\\lsass\.exe</field>
    <field name="win.eventdata.grantedAccess" type="pcre2">0x1010|0x1410|0x1fffff|0x143a</field>
    <description>CRITICAL: LSASS Memory Access - Possible Credential Dumping (Mimikatz/ProcDump)</description>
    <mitre>
      <id>T1003.001</id>
    </mitre>
    <group>credential_access,lsass,mimikatz,</group>
  </rule>

  <!-- Exclusion LSASS légitime — Sécurité AD Windows -->
  <rule id="100001b" level="0">
    <if_sid>100001</if_sid>
    <field name="win.eventdata.sourceImage" type="pcre2">(?i)(antimalware|windows\\system32\\werfault|csrss)</field>
    <description>False positive: Legitimate LSASS access from system process</description>
  </rule>

  <!-- ============================================================ -->
  <!-- T1558.003 — Kerberoasting                                    -->
  <!-- Windows Security Event ID 4769 : Kerberos Service Ticket     -->
  <!-- ============================================================ -->
  <rule id="100002" level="13">
    <if_sid>60103</if_sid>
    <field name="win.system.eventID">^4769$</field>
    <field name="win.eventdata.ticketEncryptionType">^0x17$</field>
    <description>HIGH: Kerberos RC4 Ticket Request - Possible Kerberoasting (T1558.003)</description>
    <mitre>
      <id>T1558.003</id>
    </mitre>
    <group>credential_access,kerberos,kerberoasting,</group>
  </rule>

  <!-- Kerberoasting — multiple tickets RC4 (high confidence) -->
  <rule id="100002b" level="15" frequency="5" timeframe="60">
    <if_matched_sid>100002</if_matched_sid>
    <same_field>win.eventdata.clientAddress</same_field>
    <description>CRITICAL: Multiple Kerberos RC4 requests - Active Kerberoasting Attack</description>
    <mitre>
      <id>T1558.003</id>
    </mitre>
    <group>credential_access,kerberoasting,attack_confirmed,</group>
  </rule>

  <!-- ============================================================ -->
  <!-- T1003.006 — DCSync Attack                                    -->
  <!-- Windows Security Event ID 4662 : DS Replication              -->
  <!-- ============================================================ -->
  <rule id="100003" level="15">
    <if_sid>60103</if_sid>
    <field name="win.system.eventID">^4662$</field>
    <field name="win.eventdata.properties" type="pcre2">1131f6aa-9c07-11d1-f79f-00c04fc2dcd2|1131f6ad-9c07-11d1-f79f-00c04fc2dcd2|89e95b76-444d-4c62-991a-0facbeda640c</field>
    <description>CRITICAL: DCSync Attack - DS Replication Rights Used (T1003.006)</description>
    <mitre>
      <id>T1003.006</id>
    </mitre>
    <group>credential_access,dcsync,lateral_movement,</group>
  </rule>

  <!-- ============================================================ -->
  <!-- T1550.002 — Pass-the-Hash                                    -->
  <!-- Windows Security Event ID 4624 : Logon Type 3 avec NTLM     -->
  <!-- ============================================================ -->
  <rule id="100004" level="14">
    <if_sid>60103</if_sid>
    <field name="win.system.eventID">^4624$</field>
    <field name="win.eventdata.logonType">^3$</field>
    <field name="win.eventdata.authenticationPackageName" type="pcre2">^NTLM$</field>
    <field name="win.eventdata.keyLength">^0$</field>
    <description>HIGH: Possible Pass-the-Hash - NTLM Network Logon with zero key length (T1550.002)</description>
    <mitre>
      <id>T1550.002</id>
    </mitre>
    <group>lateral_movement,pass_the_hash,ntlm,</group>
  </rule>

  <!-- ============================================================ -->
  <!-- T1218 — LOLBins Execution (Living-off-the-land)             -->
  <!-- Sysmon Event ID 1 : Process Create                           -->
  <!-- ============================================================ -->
  <rule id="100010" level="10">
    <if_group>sysmon_event1</if_group>
    <field name="win.eventdata.image" type="pcre2">(?i)(certutil\.exe|mshta\.exe|regsvr32\.exe|rundll32\.exe|wscript\.exe|cscript\.exe|installutil\.exe|regasm\.exe)</field>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)(http|https|\\\\|\.ps1|\.vbs|\.hta|download|urlcache|encode)</field>
    <description>HIGH: LOLBin Execution with suspicious arguments - Possible T1218</description>
    <mitre>
      <id>T1218</id>
    </mitre>
    <group>defense_evasion,lolbins,execution,</group>
  </rule>

  <!-- certutil download spécifique -->
  <rule id="100011" level="12">
    <if_group>sysmon_event1</if_group>
    <field name="win.eventdata.image" type="pcre2">(?i)certutil\.exe</field>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)(-urlcache|-decode|-encode|download)</field>
    <description>HIGH: CertUtil used for download/encoding - T1218.013</description>
    <mitre>
      <id>T1218.013</id>
    </mitre>
    <group>defense_evasion,certutil,lolbins,</group>
  </rule>

  <!-- ============================================================ -->
  <!-- T1021.006 — WinRM Lateral Movement                          -->
  <!-- Windows Security Event ID 4624 Logon Type 3 + WinRM port    -->
  <!-- ============================================================ -->
  <rule id="100020" level="12">
    <if_sid>60103</if_sid>
    <field name="win.system.eventID">^4624$</field>
    <field name="win.eventdata.logonType">^3$</field>
    <field name="win.eventdata.processName" type="pcre2">(?i)wsmprovhost</field>
    <description>HIGH: WinRM Session - Possible Lateral Movement (T1021.006)</description>
    <mitre>
      <id>T1021.006</id>
    </mitre>
    <group>lateral_movement,winrm,</group>
  </rule>

  <!-- ============================================================ -->
  <!-- T1053.005 — Scheduled Task / Job                            -->
  <!-- Windows Security Event ID 4698                               -->
  <!-- ============================================================ -->
  <rule id="100030" level="11">
    <if_sid>60103</if_sid>
    <field name="win.system.eventID">^4698$</field>
    <description>MEDIUM: Scheduled Task Created - Possible Persistence (T1053.005)</description>
    <mitre>
      <id>T1053.005</id>
    </mitre>
    <group>persistence,scheduled_task,</group>
  </rule>

  <!-- Scheduled task avec commande suspecte -->
  <rule id="100031" level="14">
    <if_sid>100030</if_sid>
    <field name="win.eventdata.taskContent" type="pcre2">(?i)(powershell|cmd|wscript|cscript|mshta|\\\\|http|https)</field>
    <description>HIGH: Suspicious Scheduled Task with execution payload (T1053.005)</description>
    <mitre>
      <id>T1053.005</id>
    </mitre>
    <group>persistence,scheduled_task,suspicious,</group>
  </rule>

  <!-- ============================================================ -->
  <!-- T1078.002 — Valid Accounts - Admin abuse                    -->
  <!-- Windows Security Event ID 4672 : Special Privilege Logon    -->
  <!-- ============================================================ -->
  <rule id="100040" level="9">
    <if_sid>60103</if_sid>
    <field name="win.system.eventID">^4672$</field>
    <field name="win.eventdata.privileges" type="pcre2">SeDebugPrivilege|SeTcbPrivilege|SeImpersonatePrivilege</field>
    <description>MEDIUM: High Privilege Token Assigned - Admin/Debug Privilege (T1078.002)</description>
    <mitre>
      <id>T1078.002</id>
    </mitre>
    <group>privilege_escalation,valid_accounts,</group>
  </rule>

  <!-- ============================================================ -->
  <!-- T1547.001 — Registry Run Keys (Persistence)                 -->
  <!-- Sysmon Event ID 13 : RegistryValueSet                       -->
  <!-- ============================================================ -->
  <rule id="100050" level="12">
    <if_group>sysmon_event13</if_group>
    <field name="win.eventdata.targetObject" type="pcre2">(?i)(\\CurrentVersion\\Run|\\CurrentVersion\\RunOnce|\\Windows\\CurrentVersion\\Explorer\\Shell Folders)</field>
    <field name="win.eventdata.details" type="pcre2">(?i)(\.exe|\.dll|\.ps1|\.bat|\.vbs|\.hta|cmd\.exe|powershell)</field>
    <description>HIGH: Registry Run Key Modified - Possible Persistence (T1547.001)</description>
    <mitre>
      <id>T1547.001</id>
    </mitre>
    <group>persistence,registry,</group>
  </rule>

  <!-- ============================================================ -->
  <!-- T1055 — Process Injection                                   -->
  <!-- Sysmon Event ID 8 : CreateRemoteThread                      -->
  <!-- ============================================================ -->
  <rule id="100060" level="14">
    <if_group>sysmon_event8</if_group>
    <field name="win.eventdata.targetImage" type="pcre2">(?i)(explorer\.exe|svchost\.exe|notepad\.exe|mspaint\.exe)</field>
    <description>HIGH: Remote Thread Created in System Process - Possible Process Injection (T1055)</description>
    <mitre>
      <id>T1055</id>
    </mitre>
    <group>defense_evasion,process_injection,</group>
  </rule>

  <!-- ============================================================ -->
  <!-- T1136.002 — Create Domain Account                           -->
  <!-- Windows Security Event ID 4720                               -->
  <!-- ============================================================ -->
  <rule id="100070" level="13">
    <if_sid>60103</if_sid>
    <field name="win.system.eventID">^4720$</field>
    <description>HIGH: New User Account Created in Domain - Verify Authorization (T1136.002)</description>
    <mitre>
      <id>T1136.002</id>
    </mitre>
    <group>persistence,account_manipulation,</group>
  </rule>

  <!-- ============================================================ -->
  <!-- T1562.001 — Impair Defenses: Disable Security Tools        -->
  <!-- Windows Security + PowerShell Event IDs                     -->
  <!-- ============================================================ -->
  <rule id="100080" level="14">
    <if_group>sysmon_event1</if_group>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)(Set-MpPreference|Disable-WindowsOptionalFeature|sc\.exe.*stop|net.*stop|taskkill.*defender|DisableRealtimeMonitoring)</field>
    <description>CRITICAL: Attempt to Disable Security Tools - T1562.001</description>
    <mitre>
      <id>T1562.001</id>
    </mitre>
    <group>defense_evasion,disable_security,</group>
  </rule>

</group>
